# dnsmasqdhcp/rootfs/etc/services.d/hostapd/run
#!/usr/bin/with-contenv bashio
CONFIG="/etc/hostapd.conf"

# Run dnsmasq
bashio::log.info "Starting hostapd..."

# Set max open file limit to speed up startup
ulimit -n 1024

#exec hostapd "${CONFIG}" -z < /dev/null
which hostapd
